# deploy/icloud-hme.service
[Unit]
Description=iCloud Hide My Email service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=icloud-hme
Group=icloud-hme
EnvironmentFile=/opt/icloud-hme/config/icloud-hme.env
WorkingDirectory=/opt/icloud-hme/data
ExecStart=/opt/icloud-hme/bin/icloud-hme -addr 127.0.0.1:8081 -data /opt/icloud-hme/data
Restart=on-failure
RestartSec=5s

NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/opt/icloud-hme/data

[Install]
WantedBy=multi-user.target
